# config/includes.chroot/home/user/.config/autostart/florence.desktop
[Desktop Entry]
Type=Application
Exec=florence -u /home/user/.florence.conf
