# config/includes.chroot/home/user/.config/autostart/rxvt3.desktop
[Desktop Entry]
Type=Application
Exec=rxvt -geometry 80x25+64+280 -title "help shell" -e /home/user/.helpscript
