# config/includes.chroot/home/user/.config/autostart/rxvt4.desktop
[Desktop Entry]
Type=Application
Exec=rxvt -geometry 80x25+554+280
