# config/includes.chroot/home/user/.config/autostart/rxvt2.desktop
[Desktop Entry]
Type=Application
Exec=rxvt -geometry 80x25+554+0 -title "2GiB Treasure Goblin (compatible)" -e /home/user/.tgscript 2
